# ab1_site/systemd/abi.service
[Unit]
Description=ABI Trace Viewer (Gunicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/abi/ab1_site
ExecStart=/opt/abi/ab1_site/.venv/bin/gunicorn -w 2 -b 127.0.0.1:8080 'server:create_app()'
Restart=always
RestartSec=3
User=apache
Group=apache
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

